# systemd/keycard-alert-handler.service
[Unit]
Description=Keycard Service: Dispatch alert on service failure

[Service]
Type=oneshot
EnvironmentFile=-/etc/default/keycard-service
ExecStart=/opt/keycard-service/scripts/alert %I
StandardOutput=journal
StandardError=journal

# Prevent recursive alerts: do not trigger alert handler on alert handler failure
OnFailure=
