# systemd/keycard-acl-sync.timer
[Unit]
Description=Keycard Service: Sync member ACL from Wild Apricot (every 8 minutes + daily force)

[Timer]
OnCalendar=*:0/8
Persistent=true

[Install]
WantedBy=timers.target
