# systemd/keycard-acl-sync.service
[Unit]
Description=Keycard Service: Sync member ACL from Wild Apricot
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=keycard
Group=keycard
WorkingDirectory=/opt/keycard-service
EnvironmentFile=-/etc/default/keycard-service
ExecStart=/opt/keycard-service/scripts/acl-sync
StandardOutput=journal
StandardError=journal
OnFailure=keycard-alert-handler@%n.service

[Install]
WantedBy=multi-user.target
